# packages/optilibre/optilibre-2.0.0rc2.tar.gz/optilibre-2.0.0rc2/systemd/opti_libre.timer
[Unit]
Description=opti videos

[Timer]
OnCalendar=0:45

[Install]
WantedBy=basic.target

